# docker/supervisord.ini
[supervisord]
nodaemon=true
user=root

[program:runner]
command=/bin/bash -c 'while true; do date && php /var/www/html/main.php; sleep 3600; done'
user=brokerage
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

